# install/itty-daemon.service
[Unit]
Description=iTTY Daemon
Documentation=https://github.com/honeycomb-Technologies/iTTY
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/itty start
Restart=on-failure
RestartSec=5
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# The daemon writes ~/.config/itty/config.toml and may create or update shell rc
# files via PUT /config/auto, so HOME cannot be mounted read-only here.
NoNewPrivileges=yes
ProtectSystem=strict

[Install]
WantedBy=default.target
